# deploy/lightsail/it2p.service.example
[Unit]
Description=IT2P Next.js application
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/var/www/it2p
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=HOSTNAME=127.0.0.1
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target